# audio/99-hdmi-audio.rules
ACTION=="change", SUBSYSTEM=="drm", KERNEL=="card0", ENV{HOTPLUG}=="1", RUN+="/usr/local/bin/audio-switch.sh"
